# docs/openhands/config/agent-canvas.service
[Unit]
Description=OpenHands Agent Canvas
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/spy4x
# PATH must include ~/.local/bin for uvx (agent-server runtime)
Environment=PATH=/home/spy4x/.local/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin
Environment=PORT=8001
Environment=LOCAL_BACKEND_API_KEY=REPLACE_WITH_OPENHANDS_API_KEY_FROM_SERVERS_HOME_ENV
Environment=OH_SECRET_KEY=REPLACE_WITH_OPENHANDS_SECRET_KEY_FROM_SERVERS_HOME_ENV
Environment=GITHUB_TOKEN=REPLACE_WITH_GITHUB_TOKEN_FROM_ENV_ROOT
Environment=LLM_MODEL=deepseek-chat
Environment=LLM_BASE_URL=https://api.deepseek.com/v1
Environment=LLM_API_KEY=REPLACE_WITH_DEEPSEEK_API_KEY
Environment=RUNTIME=process
Environment=AGENT_SERVER_BIND_ADDRESS=0.0.0.0

ExecStart=npx @openhands/agent-canvas --public --port 8001
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
