# ansible/templates/cloud-tunnel.service.j2
[Unit]
Description=cloud-tunnel — SSH tunnel homelab→cloudlab for mail ports
Documentation=https://github.com/spy4x/homelab
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ansible_user }}
# env -i works around OpenSSH SIGSEGV when spawned by systemd (Fedora 43 bug)
ExecStart=/usr/bin/env -i /usr/bin/ssh \
  -NT \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes \
  -o StrictHostKeyChecking=accept-new \
  -p {{ cloud_tunnel_ssh_port | default(26845) }} \
  -L 0.0.0.0:1993:{{ cloud_tunnel_target_ip }}:993 \
  -L 0.0.0.0:1587:{{ cloud_tunnel_target_ip }}:587 \
  {{ cloud_tunnel_ssh_user }}@{{ cloud_tunnel_ssh_host }}
Restart=always
RestartSec=5
StartLimitInterval=0
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
